# local.app-killer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.app-killer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/kill-apps.sh</string>
        <string>__TARGET_USER__</string>
    </array>
    <key>StartInterval</key>
    <integer>2</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/app-killer.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/app-killer.err</string>
</dict>
</plist>
